# mail-filter/pfix-srsd/files/pfix-srsd.init
#!/sbin/runscript

command="/usr/sbin/pfix-srsd"

start() {
        ebegin "Starting ${SVCNAME}"
        start-stop-daemon \
                --start \
                --exec ${command} \
                --pidfile ${PIDFILE} \
                -- \
                        -p ${PIDFILE} \
                        ${OPTS} \
                        ${DOMAIN} \
                        ${SECRETSFILE}
}
